# docker/supervisor-fishserve.conf
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n -f /etc/rsyslog.d/catchhub.conf

[program:fishserve-sender]
command = python3 manage.py run_sender --settings fishserve.settings
